# units/opt-mgw-secrets.mount
[Unit]
Before=local-fs.target

[Mount]
Type=tmpfs
What=tmpfs
Where=/opt/mgw/secrets
Options=nosuid,nodev,noexec,mode=1777,size=100M